# aconfmgr/files/etc/systemd/system/boot-zram.service
[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/system-sleep/zram-hibernate-hook post hibernate

[Install]
RequiredBy=multi-user.target
